# my-scripts/linuxtweaks-nvidia/TEST/io-scheduler.service
[Unit]
Description=Set I/O Scheduler on boot
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo mq-deadline > /sys/block/sda/queue/scheduler'

[Install]
WantedBy=multi-user.target
