# my-scripts/linuxtweaks-nvidia/TEST/linuxtweaks-cake-resume.service
# sudo systemctl daemon-reexec && sudo systemctl daemon-reload && sudo systemctl reenable linuxtweaks-cake-resume.service && sudo systemctl restart linuxtweaks-cake-resume.service && systemctl status linuxtweaks-cake-resume.service
[Unit]
Description=Tolga's V10.0 CAKE qdisc after suspend/wake for AUTO-DETECTED interface
After=network-online.target suspend.target sleep.target
Wants=network-online.target
ConditionPathExists=/usr/sbin/tc

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
StandardOutput=journal+console
StandardError=journal+console
Environment=SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=0

# Wait for interface state UP
ExecStartPre=/usr/sbin/modprobe sch_cake
ExecStartPre=/bin/bash -c 'for i in {1..10}; do interface=$(ip -o link show | awk -F": " '"'"'$2 ~ /^(wlp|wlo|wlx|eth|eno)/ && $0 ~ /state UP/ && $0 !~ /NO-CARRIER/ { print $2; exit }'"'"'); [ -n "$interface" ] && break || sleep 1; done'

# Dynamic interface detection and CAKE apply
ExecStart=/bin/bash -c 'for i in {1..5}; do interface=$(ip -o link show | awk -F": " '"'"'$2 ~ /^(wlp|wlo|wlx|eth|eno)/ && $0 ~ /state UP/ && $0 !~ /NO-CARRIER/ { print $2; exit }'"'"'); if [ -n "$interface" ]; then tc qdisc replace dev "$interface" root cake bandwidth 950Mbit besteffort diffserv4 triple-isolate nat raw ack-filter split-gso rtt 5ms overhead 48 wash && break; else sleep 2; fi; done'

# Watchdog & safety
TimeoutStartSec=10min
TimeoutStopSec=10s
TimeoutStopFailureMode=kill
SuccessExitStatus=0 3

[Install]
WantedBy=suspend.target sleep.target
