# my-scripts/linuxtweaks-nvidia/TEST/linuxtweaks-cake.service
# sudo systemctl daemon-reexec && sudo systemctl daemon-reload && sudo systemctl reenable linuxtweaks-cake.service && sudo systemctl restart linuxtweaks-cake.service && systemctl status linuxtweaks-cake.service
[Unit]
Description=Tolga's V10.0 CAKE qdisc for wlp3s0 at boot
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/sbin/tc

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
StandardOutput=journal
StandardError=journal
Environment=SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=0

# Wait for wlp3s0 state UP
ExecStartPre=/usr/sbin/modprobe sch_cake
ExecStartPre=/bin/bash -c 'for i in {1..10}; do interface=$(ip -br link show | awk '"'"'$2 == "UP" && $1 ~ /^(wlp|wlo|wlx|eth|eno)/ { print $1; exit }'"'"'); [ -n "$interface" ] && break || sleep 1; done'

# Dynamic interface detection and CAKE apply
ExecStart=/bin/bash -c 'for i in {1..5}; do interface=$(ip -br link show | awk '"'"'$2 == "UP" && $1 ~ /^(wlp|wlo|wlx|eth|eno)/ { print $1; exit }'"'"'); echo "[$(date +%%T)] Detected interface: $interface"; if [ -n "$interface" ]; then tc qdisc replace dev "$interface" root cake bandwidth 950Mbit besteffort diffserv4 triple-isolate raw nat ack-filter split-gso rtt 5ms overhead 48 wash && echo "[$(date +%%T)] CAKE applied successfully." && break; else echo "[$(date +%%T)] No interface up, retrying..."; sleep 2; fi; done'

# Watchdog & safety
TimeoutStartSec=10min
TimeoutStopSec=10s
TimeoutStopFailureMode=kill
SuccessExitStatus=0 3
Restart=on-failure

[Install]
WantedBy=multi-user.target


