# my-scripts/linuxtweaks-nvidia/TEST/my-preload.timer
[Unit]
Description=Periodic preload of apps

[Timer]
OnBootSec=2min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
